# supervisor/confluence.conf
[program:confluence]
command=/opt/atlassian/confluence/bin/start-confluence.sh -fg
priority=100
autostart=true
autorestart=true
exitcodes=0
user=confluence