# supervisor/postgres.conf
[program:postgres]
command=/usr/local/bin/docker-entrypoint.sh postgres
priority=10
autostart=true
autorestart=true
user=postgres